# scripts/install/utcoupe.conf
description	"UTCoupe server daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork

exec /root/main.sh
